# .config/autostart/polybar.desktop
[Desktop Entry]
Exec=sh ~/.config/polybar/launch.sh --forest
Icon=
Name=polybar
Path=
Terminal=False
Type=Application